# Data/user-data/User/Tools/PortableGit/etc/profile.d/aliases.sh
version https://git-lfs.github.com/spec/v1
oid sha256:a7540390bf956f33f6497a87c8986e124a8ecdacad26d0f75695c83dc2410d96
size 234
